# config/upstart/yesod-benchmark.conf
description "Website created to benchmark Yesod"

start on started mysql
stop on stopped mysql
respawn

chdir /var/www/yesod/
exec /var/www/yesod/yesod-benchmark
